# etc/rc.local
#!/bin/bash
/opt/bin/mount_workspace.sh /home/jupyter
/opt/bin/define-jupyter-service.sh
